# hass-app-testing/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bashio
exec /app/run.hass.sh
